# files/launch/devpts.conf
# @@@LICENSE
#
# Copyright (c) 2014 LG Electronics, Inc.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
# http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
# LICENSE@@@

description "Upstart version of devpts script"

start on stopped populate-volatile

console output

script
	# GID of the `tty' group
	TTYGRP=5
	# Set to 600 to have `mesg n' be the default
	TTYMODE=600

	if grep -q devpts /proc/filesystems              
	then
        	# Create multiplexor device.
        	test -c /dev/ptmx || mknod -m 666 /dev/ptmx c 5 2
                                                         
        	# Mount /dev/pts if needed.                
	        if ! grep -q devpts /proc/mounts                 
	        then                             
	                mkdir -p /dev/pts        
	                mount -t devpts devpts /dev/pts -ogid=${TTYGRP},mode=${TTYMODE}
	        fi                                                                     
	fi                       
end script
